# module/max-power-save/etc/systemd/system/power-save.timer
[Unit]
Description=Delay start of power-save by 5min

[Timer]
OnBootSec=5min
Unit=power-save.service

[Install]
WantedBy=timers.target