# module/sync-camera-data/etc/systemd/system/sync-camera-data.service
[Unit]
Description=Sync camera data

[Service]
ExecStart=/usr/local/bin/sync-camera-data

[Install]
WantedBy=multi-user.target